# systemd/jarvis-hud.service
[Unit]
Description=Jarvis HUD voice/visual command center
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/jarvis/server
EnvironmentFile=/etc/jarvis/env
ExecStart=/opt/jarvis/server/.venv/bin/python server.py
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
